# .devcontainer/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:chrome-devtools]
command=/usr/bin/google-chrome-stable --headless --disable-gpu --no-sandbox --disable-setuid-sandbox --disable-dev-shm-usage --disable-software-rasterizer --disable-extensions --remote-debugging-port=9222 about:blank
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/chrome-devtools.log
stderr_logfile=/var/log/supervisor/chrome-devtools-error.log
user=vscode
environment=HOME="/home/vscode",USER="vscode"
stopwaitsecs=10
